# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

openvpn /openvpn.conf
